# build/iso/rootfs_overlay/etc/systemd/system/ahnis.service
[Unit]
Description=Ahnis MemPalace Semantic Memory
After=network.target
[Service]
Type=simple
ExecStart=/usr/bin/python3 -m uvicorn app.main:app --host 0.0.0.0 --port 8091
WorkingDirectory=/opt/ahnis
Restart=always
User=prax
Environment=AHNIS_DB_PATH=/var/prax/ahnis.db
[Install]
WantedBy=multi-user.target
