# build/iso/rootfs_overlay/etc/systemd/system/prax.target
[Unit]
Description=Prady OS Prax AI Stack
After=multi-user.target
Wants=ollama.service neila.service ahnis.service prax-brain.service
[Install]
WantedBy=multi-user.target
