# packages/systemd/ydotoold.service
[Unit]
Description=ydotool daemon for input injection
After=basic.target

[Service]
Type=simple
ExecStart=/usr/bin/ydotoold --socket-path /tmp/.ydotool_socket \
  --socket-perm 0666
Restart=always

[Install]
WantedBy=multi-user.target
