# packages/prax-agent.service
[Unit]
Description=Prady AI Agent Orchestration Service
After=network.target vyrex.service prady-input.service
Requires=vyrex.service

[Service]
Type=notify
User=prady
Group=prady
WorkingDirectory=/opt/prady/agent
ExecStart=/usr/bin/node /opt/prady/agent/dist/index.js
Restart=on-failure
RestartSec=5
Environment=NODE_ENV=production
Environment=VYREX_SOCKET=/run/vyrex/api.sock
Environment=PRADY_DBUS=org.prady.Input
Environment=PRADY_MEMORY_DB=/var/prady/memory.db
PrivateTmp=true
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/var/prady /tmp

[Install]
WantedBy=multi-user.target
